# services/hdd-temp.timer
[Unit]
Description="硬盘温度监控"

[Install]
WantedBy=timers.target

[Timer]
OnBootSec=1min
OnUnitInactiveSec=30s
OnUnitActiveSec=30s
